# bluetooth/telem-server.service
[Unit]
Description=Race telemetry server
After=network.target

[Service]
Type=simple
User=gearados
WorkingDirectory=/home/gearados/repos/telem/server
ExecStart=/bin/bash -c 'source /home/gearados/.nvm/nvm.sh && exec npx tsx src/main.ts'
Environment=DATA_DIR=/home/gearados/repos/telem/server/data
Environment=PORT=4400
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
